# bac.service
[Unit]
Description=Bac

[Service]
ExecStart=PATH_TO_/bac.sh